# scripts/argusnexus-survivor-paper.service
[Unit]
Description=ArgusNexus V4 Paper Trader with Learning System
Documentation=http://localhost:8000
After=network.target

[Service]
Type=simple
User=your_user
Group=tony
WorkingDirectory=/path/to/ArgusNexus
ExecStart=/path/to/ArgusNexus/venv/bin/python /path/to/ArgusNexus/scripts/live_paper_trader.py --interval 1h --capital 500
Restart=always
RestartSec=30

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/path/to/ArgusNexus

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argusnexus-paper

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/path/to/ArgusNexus/data

# Timeouts
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
