# scripts/argusnexus-v4-paper.service
[Unit]
Description=ArgusNexus V4 Paper Trader with Learning System
After=network.target

[Service]
Type=simple
User=your_user
WorkingDirectory=/path/to/ArgusNexus
ExecStart=/path/to/ArgusNexus/venv/bin/python /path/to/ArgusNexus/scripts/live_paper_trader.py --interval 4h --capital 10000
Restart=always
RestartSec=30

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/path/to/ArgusNexus

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argusnexus-paper

# Timeouts
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target
