# scripts/argusnexus-dashboard-api.service
[Unit]
Description=ArgusNexus V4 Dashboard API - Bloomberg Terminal
Documentation=http://localhost:8000
After=network.target

[Service]
Type=simple
User=your_user
Group=tony
WorkingDirectory=/path/to/ArgusNexus
ExecStart=/path/to/ArgusNexus/venv/bin/uvicorn src.api.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/path/to/ArgusNexus

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argusnexus-api

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/path/to/ArgusNexus/data

# Timeouts
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
